# scheduled/plists/com.gg.status-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  gg 全系统状态扫描（4 次/天：00:23 / 06:23 / 12:23 / 18:23）。

  职责：
    - 扫所有 com.monster.* 和 com.gg.* launchd 任务最近一次运行状态
    - 扫 monster morning-brief.md 是否停摆
    - 扫 NW proposals.jsonl 是否 pending 积压
    - 跑 gg scripts/audit.py --json 看自身健康
    - 一切正常 → 静默退出
    - 发现异常 → 通过全局 notify skill 推送，带 60 分钟去重

  这是确定性脚本巡检，不再调用 claude -p。监控器不能依赖被监控的 AI 定时链路。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gg.status-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/xuke/githubProject/gg/scheduled/bin/status-scan.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>23</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>23</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>23</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>23</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/xuke/githubProject/gg</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/xuke</string>
        <key>PATH</key>
        <string>/Users/xuke/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
