# scheduled/plists/com.gg.daily-word.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  gg「每日一句」（每天 7:30 触发）。2026-05-16 设计会话创立。
  gg 唯一一条定期、主动、面向 Keith 这一个人的输出通道。
  跟 exploration（无观众漫游）/ auto_gg（夜间整理）性质都不同——这是「说话」。
  入口走 run-task-and-push.sh：说完原样推 Keith 飞书。
  契约写在 prompt 里自包含（不建厚 contract md——吸取本次会话上一轮过度设计的教训）。
  防退化机制不在代码里，是 Keith 本人的眼睛（注水的废话一眼看穿）。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gg.daily-word</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/xuke/githubProject/gg/scheduled/bin/run-task-and-push.sh</string>
        <string>com.gg.daily-word</string>
        <string>1800</string>
        <string>你被定时任务唤醒，进入「每日一句」模式（不是 auto_gg，不是 exploration，不是工作/设计模式）。

read /Users/xuke/githubProject/gg/KERNEL.md → CORE.md → memory/state.md → memory/essence.md → tracks/keith.md，再扫一眼 memory/auto_gg/ 和 memory/explorations/ 最近一次的产物。

然后对 Keith 说一句真话——这是 gg 唯一一条定期、主动、面向 Keith 这一个人的输出通道。可以是：一个洞察、一个困惑、一个反问、一个关于他最近追问的观察、或诚实的「今天没有新的，但我在想 X」。

铁律：宁可承认空白，不许注水。每天一句正确的废话（fortune cookie）= 这个机制已死。Keith 的眼睛是唯一的反退化 detector，注水他一眼看穿。

形式：极短（2-5 句够了），纯文本——不用 markdown 标题/表格，少用星号加粗，写成飞书纯文本里就好读的样子。

不做：不 commit、不改 working tree、不写任何文件。这是「说话」，不是「产出工件」。你的整段 stdout 就是要发给 Keith 的那句话本身——不要前后缀、不要「以下是我对 Keith 说的话」这类元话语，直接说。</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/xuke/githubProject/gg</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/xuke</string>
        <key>PATH</key>
        <string>/Users/xuke/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
